# vhost.conf
<VirtualHost *:80>
    ServerAdmin webmaster@cronlms.tk
    ServerName cronlms.tk
    DocumentRoot /var/www/cronlms/current/public
    <Directory /var/www/cronlms/current/public">
        Options Indexes FollowSymLinks
        AllowOverride All
        Order allow,deny
        Allow from all
    </Directory>
</VirtualHost>

[program:lmsprod-websockets-listener-prod]
command=php /var/www/enovusptop_prod/current/artisan websockets:serve --port=6020

[program:lmsprod-queue-listener-prod]
command=php /var/www/enovusptop_prod/current/artisan queue:work --tries=5 --env=production --timeout=0


php artisan cronmon:createadmin jcolmenares jcolmenares@enovus.cl